# recipes-custom/ifupdown/files/networking.service
[Unit]
Description=ifup network start

[Service]
Type=notify
ExecStart=/sbin/ifup -a

[Install]
WantedBy=multi-user.target
