# recipes-images/base-files/base-files/cam.service
[Unit]
Description=oscam

[Service]
Type=simple
ExecStart=-/usr/bin/oscam
ExecReload=/bin/kill -HUP $MAINPID
 
[Install]
WantedBy=multi-user.target
