# recipes-images/base-files/base-files/fstrim.service
[Unit]
Description=Runs fstrim on all mounted devices that support TRIM

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/sbin/fstrim -a'

