# recipes-images/base-files/base-files/lcd.service
[Unit]
Description=configure lcd 
Requires=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lcd.sh

[Install]
WantedBy=multi-user.target

