# recipes-images/base-files/base-files/media-HDD.mount
[Unit]
Description=Mount Root Partition
Before=local-fs.target

[Mount]
Where=/media/HDD
What=/dev/sda1

[Install]
WantedBy=multi-user.target

