# recipes-images/base-files/base-files/mnt-partition_4.automount
[Unit]
Description=Mount Image Partitions 

[Automount]
Where=/mnt/partition_4

[Install]
WantedBy=multi-user.target

