# recipes-images/base-files/base-files/mnt-partition_2.automount
[Unit]
Description=Mount Image Partitions 
[Automount]
Where=/mnt/partition_2

[Install]
WantedBy=multi-user.target


