# recipes-images/base-files/base-files/mnt-partition_3.automount
[Unit]
Description=Mount Image Partitions 

[Automount]
Where=/mnt/partition_3

[Install]
WantedBy=multi-user.target


