# recipes-images/base-files/base-files/mnt-partition_1.automount
[Unit]
Description=Mount Image Partitions

[Automount]
Where=/mnt/partition_1

[Install]
WantedBy=multi-user.target


